# board/common/rootfs/etc/rc.local
#!/bin/sh
# Called after all other services have been started by Finit

# Set up firewall rules, port forwards, etc.
#nft -f /etc/nftables.conf

exit 0